# build/pabx/supervisor.conf
[supervisord]
nodaemon=true
loglevel=debug
pidfile=/tmp/supervisord.pid
logfile=/var/log/supervisord.log

[include]
files = /etc/supervisor/conf.d/*.conf

[program:php-fpm7.0]
command=php-fpm
autostart=true
autorestart=true
startretries=1
startsecs=1
redirect_stderr=true
user=root
killasgroup=true
stopasgroup=true
priority=5

[program:asterisk]
command=asterisk
autostart=true
autorestart=true
startretries=1
startsecs=1
redirect_stderr=true
user=root
killasgroup=true
stopasgroup=true
priority=10
